# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/config/logs/supervisord.log
pidfile=/config/supervisord.pid

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:django]
command=gunicorn --bind 127.0.0.1:8001 --workers 3 --timeout 120 --access-logfile /config/logs/gunicorn-access.log --error-logfile /config/logs/gunicorn-error.log suggesterr.wsgi:application
directory=/app
user=suggesterr
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/config/logs/django.log